# roles/web/files/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user=ubuntu                 ; (default is current user, required if root)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:hikerpix]
command=/home/ubuntu/hikerpix/venv/bin/uwsgi
  --virtualenv /home/ubuntu/hikerpix/venv
  --module wsgihandler
  --master
  --pidfile /tmp/hikerpix-master.pid
  --touch-reload /home/ubuntu/hikerpix/web2py/applications/hikerpix/reload
  --processes 2
  --socket 127.0.0.1:3031
  --stats 127.0.0.1:3032
  --harakiri 60
  --harakiri-verbose
  --reload-mercy 10
  --reload-on-as 256
  --reload-on-rss 192
  --vacuum
directory=/home/ubuntu/hikerpix/web2py
stopsignal=QUIT

[program:celery]
environment=PATH="/home/ubuntu/hikerpix/venv/bin"
command=nice -n 10 /home/ubuntu/hikerpix/venv/bin/celeryd
directory=/home/ubuntu/hikerpix/web2py/applications/hikerpix/conf/celery
